# home/.config/systemd/user/battery.service
[Unit]
Description = Check battery level.

[Service]
Type = oneshot
ExecStart = /home/%u/bin/k-check-battery
